# etc/s6/envoy/run
#!/bin/sh

set -e

/usr/local/bin/envoy -c /etc/envoy/envoy.yaml
